# autostart/nvidia-preferred-mode.desktop
[Desktop Entry]
Type=Application
Name=nvidia=preferred-mode
Exec=sh /home/pomp/.nvidia-preferred-mode.sh
StartupNotify=false
Terminal=false
