# home/.config/autostart/activitywatch.desktop
[Desktop Entry]
Type=Application
Encoding=UTF-8
Name=ActivityWatch
Comment=ActivityWatch
Terminal=false
Categories=Utility;
Exec=aw-qt
