# home/.config/autostart/pomky.desktop
[Desktop Entry]
Type=Application
Name=pomky
Exec=/home/pomp/.cargo/bin/pomky
StartupNotify=false
Terminal=false
Categories=System;Monitor;
