# home/.config/autostart/pomky.desktop
[Desktop Entry]
Type=Application
Name=pomky
Exec=pomky
StartupNotify=false
Terminal=false
Categories=System;Monitor;
